# infra/com.nousergon.symposion-autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nousergon.symposion-autoupdate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/brianmcmahon/Development/symposion/infra/auto-update.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/brianmcmahon/Development/symposion</string>

  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/brianmcmahon/Library/Logs/symposion-autoupdate.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/brianmcmahon/Library/Logs/symposion-autoupdate.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/brianmcmahon/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
